# config/crontab
# Crontab for the review-app `cron` server role (config/deploy.review.yml).
#
# Loaded into the cron container via:
#   (env && cat config/crontab) | crontab - && cron -f
#
# The `env` prefix copies the container's env vars (RAILS_ENV, the Postgres /
# Redis addresses, secrets) into the crontab so jobs inherit them — cron
# otherwise runs jobs with a bare environment. Each job still cd's into /rails
# for the right working directory, and redirects output to PID 1's stdout/stderr
# so it shows up in `kamal logs`.
#
# This file is only used by Kamal review apps, never by production (Cloud66).
#
# `flock -n` guards each job: every run is a fresh `bundle exec rake` (a full
# ~115MB Rails boot), so if one runs long (e.g. the host is busy) the next
# minute's run would start before it finishes and they pile up — enough stacked
# boots OOM the host. The non-blocking lock makes an overlapping run exit
# immediately instead.
*/1  * * * * cd /rails && flock -n /tmp/run_scheduler.lock        bundle exec rake run_scheduler              > /proc/1/fd/1 2>&1
*/5  * * * * cd /rails && flock -n /tmp/pghero.lock               bundle exec rake pghero:capture_query_stats > /proc/1/fd/1 2>&1
*/30 * * * * cd /rails && flock -n /tmp/read_logged_searches.lock bundle exec rake read_logged_searches       > /proc/1/fd/1 2>&1
